# browser-extension/browser-extension-bridge.service
[Unit]
Description=Claude Code Browser Extension Bridge Server
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/edu/autorenta/browser-extension
ExecStart=/usr/bin/node /home/edu/autorenta/browser-extension/bridge-server.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
