# apps/web/tools/archive_p2p/systemd/p2p-executor.service
[Unit]
Description=P2P Automation - MercadoPago Executor
After=network.target
Documentation=https://github.com/autorenta/p2p-automation

[Service]
Type=simple
User=edu
Group=edu
WorkingDirectory=/home/edu/autorenta/apps/web/tools/p2p-automation
Environment=NODE_ENV=production
Environment=DISPLAY=:0
ExecStart=/usr/bin/node dist/index.js executor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
